# ansible-deploy/roles/nodejs_app/templates/payment-gateway.service.j2
[Unit]
Description=Payment Gateway Simulator (Node.js / Express)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ pg_run_user }}
Group={{ pg_run_group }}
WorkingDirectory={{ pg_app_dir }}
Environment=NODE_ENV=production
Environment=PG_HOST={{ pg_host }}
Environment=PG_PORT={{ pg_port }}
Environment=PG_SUCCESS_RATE={{ pg_success_rate }}
Environment=PG_APP_LOG={{ pg_log_dir }}/payment-app.log
Environment=PG_ACCESS_LOG={{ pg_log_dir }}/access.log
ExecStart=/usr/bin/node {{ pg_app_dir }}/src/app.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
